# infra/Seguridad/guemes/services/clamav/clamav-log-cleaner.service
[Unit]
Description=Limpieza de logs viejos de ClamAV

[Service]
Type=oneshot
ExecStart=/usr/local/bin/clamav-log-cleaner
